# scripts/ns-core-nic.service
[Unit]
Description=Script to manage Netronome Core NIC firmware and udev rules

[Service]
ExecStart=/opt/netronome/bin/ns-core-nic-setup boot
RemainAfterExit=yes
Type=oneshot

[Install]
WantedBy=multi-user.target
